# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/usr/bin/python3 /app/Quetzal/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
